# run/systemd/generator/boot.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
Where=/boot
What=/dev/disk/by-uuid/42431555-cf82-4517-a1dc-9c884511c1ca
Type=xfs
